# tools/com.photovideo.refresh-cookies.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.photovideo.refresh-cookies</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/moykin/Library/Scripts/photovideo-refresh-cookies.sh</string>
    </array>

    <!-- Run every 24 hours -->
    <key>StartInterval</key>
    <integer>86400</integer>

    <!-- Also run immediately when loaded -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/photovideo-cookie-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/photovideo-cookie-refresh.log</string>
</dict>
</plist>
